# deploy/paddle-ocr-api.service
[Unit]
Description=PaddleOCR API — OCR 文字识别服务
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data

# 项目路径 — 按实际部署路径修改
WorkingDirectory=/opt/paddle-ocr-api
Environment=PATH=/opt/paddle-ocr-api/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin

# 启动命令
ExecStart=/opt/paddle-ocr-api/.venv/bin/uvicorn paddle_ocr_api.app:app \
    --host 0.0.0.0 \
    --port 8100 \
    --workers 1 \
    --log-level info

# 重启策略
Restart=always
RestartSec=5

# 日志
StandardOutput=journal
StandardError=journal

# 安全
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
